# telemetry/thingsboard_telemetry.service
[Unit]
Description=thingsboard_telemetry services.
Requires=docker.service
After=docker.service

[Service]
WorkingDirectory=/home/pi/chirpstack-docker/telemetry
ExecStart=/bin/bash telemetry.sh
RestartSec=3
Restart=always


[Install]
WantedBy=multi-user.target
